# etc/profile.d/sudo-doas.sh
alias sudo=sudo-doas
alias doas=sudo-doas
